# scripts/setup-tiny.d/eth-eni-setup
#!/sbin/openrc-run

# vim: set ts=4 et:
# shellcheck shell=sh

description="Sets up interfaces for attached Elastic Network Interfaces"

depend() {
    before net
    need sysfs
}

start() {
    local iface

    ebegin "Setting up interfaces for attached ENIs"
    /etc/network/assemble-interfaces
    eend "$?"
}
