# scripts/setup-ami.d/eth-eni-setup
#!/sbin/openrc-run

description="Sets up interfaces for attached Elastic Network Interfaces"

depend() {
	before net
	need sysfs
}

start() {
	local iface

	ebegin "Setting up interfaces for attached ENIs"
  /etc/network/assemble-interfaces
	eend "$?"
}
